# contrib/upstart/governator.conf
# Governator Service

description     "Governator Process Manager"
author          "Alberto G. Hierro <alberto@garciahierro.com>"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 2 5

kill timeout 300

exec /usr/local/bin/governator -D
